# lepton-view.service
[Unit]
Description=Lepton Thermal Provider (Qt linuxfb)
After=local-fs.target systemd-udev-trigger.service
Wants=systemd-udev-trigger.service

[Service]
Type=simple
# Adjust path if user cloned elsewhere, but we assume ~/RESOFLY
WorkingDirectory=/home/team13/RESOFLY/Raspberry-Pi-FLIR-Lepton-Thermal-Imaging-Camera/LeptonModule/software/raspberrypi_video

# Run on Framebuffer (Headless-ish)
Environment=QT_QPA_PLATFORM=linuxfb
Environment=QT_QPA_FB=/dev/fb0
Environment=XDG_RUNTIME_DIR=/tmp/runtime-root

# Clean/Create FIFO for commands
ExecStartPre=/bin/rm -f /tmp/lepton_cmd
ExecStartPre=/usr/bin/mkfifo /tmp/lepton_cmd
ExecStartPre=/bin/chmod 666 /tmp/lepton_cmd

# Start the App (Starts MJPEG on 8080)
ExecStart=/home/team13/RESOFLY/Raspberry-Pi-FLIR-Lepton-Thermal-Imaging-Camera/LeptonModule/software/raspberrypi_video/raspberrypi_video -cm 3 -ss 20

Restart=always
RestartSec=2
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
