# resofly.service
[Unit]
Description=RESOFLY Service
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/thermo-vision-hub/backend
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=10
Environment=PORT=8000
# Add any other env vars here, e.g.:
# Environment=DATABASE_URL=sqlite+aiosqlite:///...

[Install]
WantedBy=multi-user.target
