# resofly-tunnel.service
[Unit]
Description=ResoFly Cloudflare Tunnel & Emailer
After=network.target resofly.service

[Service]
User=team13
WorkingDirectory=/home/team13/RESOFLY/thermo-vision-hub/backend
ExecStart=/bin/bash start_tunnel_and_email.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
